# udev/10-usbkey.rules
ACTION=="add", SUBSYSTEM=="block", RUN+="usbkey unlock --device $devnode"
