# udev/10-usbkey.rules
ACTION=="add", SUBSYSTEM=="block", RUN+="/usr/local/sbin/usbkey unlock --device $devnode"
